# deploy/portfoilos-backend.service
[Unit]
Description=Portfoilos Django Backend (Gunicorn)
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/home/portfoilos/portfoilos/portfoilo_tenant_server

# Load the production .env file
EnvironmentFile=/home/portfoilos/portfoilos/portfoilo_tenant_server/.env

# Use the production settings module
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"

# Gunicorn: 3 workers (2*CPU + 1 is the rule of thumb for a 1-CPU VPS)
ExecStart=/home/portfoilos/portfoilos/portfoilo_tenant_server/.venv/bin/gunicorn \
    --workers 3 \
    --worker-class sync \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/portfoilos/gunicorn-access.log \
    --error-logfile  /var/log/portfoilos/gunicorn-error.log \
    --log-level warning \
    config.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
